# usb-notify.service
[Unit]
Description=USB notify service

[Service]
ExecStart=/usr/bin/usb-notify

[Install]
WantedBy=graphical.target
